# deploy/iva-reminders.timer
[Unit]
Description=Iva reminders dispatcher — every 5 minutes

[Timer]
OnBootSec=30s
OnCalendar=*-*-* *:0/5:00
AccuracySec=30s
Persistent=true
Unit=iva-reminders.service

[Install]
WantedBy=timers.target
